# systemd/telegram-parilka-mcp-sync.service
[Unit]
Description=Telegram Parilka MCP SQLite cache warmer
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=10min
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory=/root/telegram-parilka-mcp
ExecStartPre=/root/telegram-parilka-mcp/bin/telegram-parilka-mcp-check-build /root/telegram-parilka-mcp /root/telegram-parilka-mcp/dist/sync-daemon.js
ExecStart=/root/telegram-parilka-mcp/bin/telegram-parilka-mcp-sync-daemon
Restart=on-failure
RestartSec=30
Environment=HOME=/root

[Install]
WantedBy=default.target
